# dashlive.service
[Unit]
Description=dashlive

[Service]
Type=simple
ExecStart=/home/m/dash-live-source-simulator/run_service.sh

[Install]
WantedBy=multi-user.target